# deploy/com.aaronsun.football-predictor-claude.backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Daily results-backfill for football-predictor-claude.

        Why: the in-process APScheduler runs `backfill_recent_results` daily,
        but only while the server is up. When the server has been down for a
        few days, the matches table lags behind reality. Without this,
        prediction-audit can't reconcile yesterday's fixtures, and the audit
        tab on the web UI shows zero new resolved predictions for days.

        Schedule: 06:00 local. football-data.co.uk usually has weekend results
        within hours, but ClubElo refreshes early-morning UTC, so 06:00 local
        catches both.

        Install:
            cp deploy/com.aaronsun.football-predictor-claude.backfill.plist \
               ~/Library/LaunchAgents/
            launchctl load ~/Library/LaunchAgents/com.aaronsun.football-predictor-claude.backfill.plist

        Verify:
            tail -f /tmp/football-predictor-claude.backfill.log
    -->
    <key>Label</key>
    <string>com.aaronsun.football-predictor-claude.backfill</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/seek/Documents/football-predictor-claude/.venv/bin/python</string>
        <string>predict.py</string>
        <string>backfill-results</string>
        <string>--days</string><string>7</string>
        <string>--quiet</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/seek/Documents/football-predictor-claude</string>

    <!-- 06:00 local time, once a day. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>6</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <!-- Don't fire on load — backfill is heavier than snapshot and we don't
         want a thundering herd if both plists are installed simultaneously. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/football-predictor-claude.backfill.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/football-predictor-claude.backfill.log</string>

    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>
